# collector-health.service
[Unit]
Description=Check collector health

[Service]
Type=oneshot
ExecStart=/root/kalshi-trader/collector-health.sh
