# trader.service
[Unit]
Description=Kalshi WeatherMaker Trader
After=network-online.target collector.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/kalshi-trader
EnvironmentFile=/root/kalshi-trader/.env
ExecStart=/root/kalshi-trader/.venv/bin/python scripts/run_trader.py --environment demo
Restart=on-failure
RestartSec=30
StandardOutput=append:/root/kalshi-trader/trader.log
StandardError=append:/root/kalshi-trader/trader.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
